# docs/deploy/stratum-bot.service
# systemd "user" service template for the stratum Q&A Discord bot (ADR-0010).
# Do NOT edit this file directly on the VM — deploy.sh renders the __PLACEHOLDERS__
# and installs the result to ~/.config/systemd/user/stratum-bot.service.
#
# 常駐は社内 Ubuntu VM の systemd user サービス(no-root / no-Docker)。
# Restart=always でクラッシュ自動復帰、WantedBy=default.target + `loginctl enable-linger`
# でログオフ/再起動後も常駐する(ADR-0010 D1/D4)。
[Unit]
Description=stratum Q&A Discord bot
Documentation=https://github.com/queeenb-com/knowledge-platform/blob/main/docs/deploy/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__REPO_DIR__
# シークレット(Claude on AWS / Discord)は chmod 600 の env ファイルから読む(§9.1 / ADR-0010 D2)。
EnvironmentFile=__ENV_FILE__
# 非シークレットのパスはここで明示(env ファイルには秘密だけを置く)。
Environment=CONFIG_DIR=__REPO_DIR__/apps/discord-bot/config
Environment=CLONES_DIR=__DATA_DIR__/clones
Environment=DB_PATH=__DATA_DIR__/data/bot.db
Environment=PROMPTS_DIR=__REPO_DIR__/prompts
ExecStart=__NODE_BIN__ apps/discord-bot/dist/index.js
Restart=always
RestartSec=5
# ログは journald(`journalctl --user -u stratum-bot`)。
SyslogIdentifier=stratum-bot

[Install]
WantedBy=default.target
